# scripts/verify-prayer-times/au.altaqwa.prayer-times-verify.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>au.altaqwa.prayer-times-verify</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/shaqirfarook/prayer-times-app/scripts/verify-prayer-times/.venv/bin/python</string>
        <string>/Users/shaqirfarook/prayer-times-app/scripts/verify-prayer-times/verify.py</string>
        <string>--push</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/shaqirfarook/prayer-times-app/scripts/verify-prayer-times/verify.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/shaqirfarook/prayer-times-app/scripts/verify-prayer-times/verify.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
